# cron.d/postgresql
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin


# Sauvegardes PostgreSQL
# ======================

0 0 * * * 	root	/usr/local/bin/backup-postgresql.sh yousurf
5 0 * * * 	root	/usr/local/bin/backup-postgresql.sh smarthome


